# build/goupile/server/goupile/files/goupile@.service
[Service]
Type = simple

User = goupile
Group = goupile

RuntimeDirectory = goupile
ExecStart = /opt/goupile/domains/%i/goupile
WorkingDirectory = /opt/goupile/domains/%i

Restart = on-failure
TimeoutStopSec = 30

[Install]
WantedBy = multi-user.target
