# build/goupile/goupile@.service
[Service]
Type=simple

User=goupile
Group=goupile

ExecStart=/srv/www/goupile/bin/goupile -C %i/goupile.ini --migrate
WorkingDirectory=/srv/www/goupile/instances

Restart=on-failure
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
