# deploy/ansible/roles/borg/templates/borgmatic.timer.j2
[Unit]
Description = Run Borgmatic backup

[Timer]
OnCalendar = {{ borg_time }}
Persistent = true

[Install]
WantedBy = timers.target
