# build/goupile/server/goupile@.service
[Service]
Type=simple

User=goupile
Group=goupile

PermissionsStartOnly=true
ExecStartPre=-/bin/mkdir -p -m 0755 /run/goupile
ExecStartPre=-/bin/chown goupile:goupile /run/goupile
ExecStart=/srv/www/goupile/bin/goupile -C %i/goupile.ini --migrate
WorkingDirectory=/srv/www/goupile/domains

Restart=on-failure
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
