# build/goupile/server/goupile/templates/goupile@.service.j2
[Service]
Type = simple

User = {{goupile_run_user}}

RuntimeDirectory = goupile
ExecStart = {{goupile_root_path}}/domains/%i/goupile
WorkingDirectory = {{goupile_root_path}}/domains/%i

Restart = on-failure
TimeoutStopSec = 30

[Install]
WantedBy = multi-user.target
