# build/goupile/server/goupile/templates/goupile@.service.j2
[Unit]
Description = Goupile eCRF (%i)
Documentation = https://goupile.fr/

[Service]
Type = simple
User = {{goupile_run_user}}
Group = {{goupile_run_group}}
RuntimeDirectory = goupile
ExecStart = /usr/sbin/runc --root /run/goupile run runc_%i
WorkingDirectory = {{goupile_root_path}}/bundles/%i
Restart = on-failure
TimeoutStopSec = 30
LimitNOFILE = 16384

# Allow rootless runC to manipulate namespaces despite
# kernel.unprivileged_userns_clone=0.
# runC will then run goupile without any capability.
SecureBits = keep-caps
CapabilityBoundingSet = CAP_SYS_ADMIN
AmbientCapabilities = CAP_SYS_ADMIN

[Install]
WantedBy = multi-user.target
