# src/meestic/dist/linux/meestic.service
[Unit]
Description = Meestic daemon
Documentation = https://koromix.dev/

[Service]
Type = notify
ExecStart = /usr/bin/meestic daemon -C /etc/meestic.ini --sandbox
Restart = on-failure
RestartSec = 20

[Install]
WantedBy = multi-user.target
