# src/rekkord/dist/debian/RekkordDrop.service
[Unit]
Description=Rekkord Drop daemon
Documentation=https://rekkord.org/

[Service]
Type=notify
User=rekkorddrop
Group=rekkorddrop
ExecStart=/usr/bin/RekkordDrop -C /etc/rekkord/RekkordDrop.ini --sandbox
Restart=on-failure
RestartSec=20

[Install]
WantedBy=multi-user.target
