# src/heimdall/dist/linux/heimdall.service
[Unit]
Description=Heimdall daemon
Documentation=https://koromix.dev/

[Service]
Type=notify
User=heimdall
Group=heimdall
UMask=0002
ExecStart=/usr/bin/heimdall -C /etc/heimdall.ini --sandbox
Restart=on-failure
RestartSec=20

[Install]
WantedBy = multi-user.target
